# src/embed/gitsitter.service
# Keep in sync with nix/home-manager-module.nix (systemd.user.services.gitsitter)
[Unit]
Description=gitsitter daemon

[Service]
Type=exec
ExecStart=@@EXEC_PATH@@ daemon run
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
